# pod/primary/supervisord.conf
; pod/primary/supervisord.conf — Phase 6.6 Wave 0 LOCKED.
;
; Replaces the rejected DinD design (06.6-SPIKE-dind-privileged.md proved
; overlayfs in nested namespace fails on Vast.ai) with supervisord as
; PID 1 + 4 co-located child processes (06.6-WAVE0-GATES.md Decision 2).
; Children share the container's network namespace (localhost ports),
; the GPU device, and the filesystem — same intent as Strategy B
; (4 services in 1 pod for low intra-pod latency + GPU sharing).
;
; The jinja flag on [program:llama] is the B1 GGUF-embedded path locked
; by 06.6-WAVE0-GATES.md Decision 3 + 06.6-SPIKE-qwen3.6-jinja.md Round 3
; (Qwen3.6 GGUF embedded chat_template renders tool-calling end-to-end
; via the jinja flag alone; no external template path needed).
;
; This file is baked into the image at build time (pod/primary/Dockerfile
; COPY directive). Runtime overrides are NOT supported in this plan —
; that would require regenerating supervisord.conf inside the pod and
; restarting supervisord, adding fragility. Deferred per CONTEXT.md
; Deferred Ideas to a future B2 fallback path.

[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:llama]
command=/app/llama-server --host 0.0.0.0 --port 8000 -m /weights/qwen/model.gguf -ngl 99 -np 2 --ctx-size 16384 --jinja
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/llama.log
priority=100

[program:speaches]
command=/opt/speaches-venv/bin/uvicorn speaches.main:create_app --host 0.0.0.0 --port 8001 --factory
environment=WHISPER__MODEL_DIR="/weights/whisper",SPEACHES__MODEL_DIR="/weights/whisper"
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/speaches.log
priority=200

[program:infinity]
command=/opt/infinity-venv/bin/infinity_emb v2 --port 8002 --model-id /weights/bge-m3
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/infinity.log
priority=300

[program:dcgm]
command=/usr/bin/dcgm-exporter --address 0.0.0.0:9400
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/dcgm.log
priority=400
