# pod/primary/supervisord.conf
; pod/primary/supervisord.conf — Phase 11.2 RESTORED (tier-0 local STT back).
;
; Phase 11.2 (2026-06-07) RESTORES [program:speaches] (faster-whisper tier-0
; STT) verbatim from commit 0f65da5^ — the Phase 11.1 removal is reverted
; because the tier-1 OpenAI-only fallback proved too costly ($0.36/h). The
; new tier-1 cascade (gemini-stt → groq-whisper → openai-whisper) is wired
; in the gateway (Phase 11.2 Plan 03) and tier-0 local-stt is preferred
; whenever the primary pod is ON. Child set: llama + speaches + chatterbox +
; dcgm (health-bridge runs separately in the pod, not under supervisord).
; The [program:infinity] commentary block that remained post-Phase 06.7 D-03
; is gone (D-A2 bundle — embed stays off-pod).
;
; Replaces the rejected DinD design (06.6-SPIKE-dind-privileged.md proved
; overlayfs in nested namespace fails on Vast.ai) with supervisord as
; PID 1 + co-located child processes (06.6-WAVE0-GATES.md Decision 2).
; Children share the container's network namespace (localhost ports),
; the GPU device, and the filesystem.
;
; The jinja flag on [program:llama] is the B1 GGUF-embedded path locked
; by 06.6-WAVE0-GATES.md Decision 3 + 06.6-SPIKE-qwen3.6-jinja.md Round 3
; (Qwen3.6 GGUF embedded chat_template renders tool-calling end-to-end
; via the jinja flag alone; no external template path needed).
;
; This file is baked into the image at build time (pod/primary/Dockerfile
; COPY directive). Runtime overrides are NOT supported in this plan —
; that would require regenerating supervisord.conf inside the pod and
; restarting supervisord, adding fragility. Deferred per CONTEXT.md
; Deferred Ideas to a future B2 fallback path.

[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:llama]
command=/app/llama-server --host 0.0.0.0 --port 8000 -m /weights/qwen/model.gguf -ngl 99 -np 2 --ctx-size 16384 --jinja
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/llama.log
priority=100

[program:speaches]
command=/opt/speaches-venv/bin/uvicorn speaches.main:create_app --host 0.0.0.0 --port 8001 --factory
environment=HF_HUB_CACHE="/weights/whisper"
; Working directory must contain `realtime-console/dist` (relative path that
; speaches 0.9.0rc3 mounts via StaticFiles at /v1/realtime — see Dockerfile
; mkdir of /opt/speaches-data/realtime-console/dist + stub index.html).
; Boots with `RuntimeError: Directory 'realtime-console/dist' does not exist`
; if directory= is missing (UAT 2026-05-18 attempt 11).
directory=/opt/speaches-data
; Phase 14 SEED-019 part 2 (2026-06-17): the hardcoded
; WHISPER__INFERENCE_DEVICE="cpu" pin was REMOVED from the environment= line
; above. The whisper device is now VRAM-adaptive: the gateway-built onstart
; (gateway/internal/primary/onstart.go) reads total VRAM via nvidia-smi and
; exports WHISPER__INFERENCE_DEVICE (cuda on >=30 GB shapes pinned to the
; max-free GPU via WHISPER__DEVICE_INDEX, cpu below threshold) BEFORE
; `exec supervisord`. With the key absent from environment= here,
; [program:speaches] INHERITS the onstart-exported value (supervisord env
; inheritance — Option A). If the var were ever unset, speaches fails OPEN to
; its own defaults inference_device=auto / device_index=0 (never crashes).
; History: UAT 17 (2026-05-19) added the cpu pin on RTX 4090 24 GB (no headroom
; for whisper GPU offload alongside Qwen + bge-m3 + KV — UAT 16 CUDA OOM);
; Phase 11.2 kept it for the 2x3090 default. Phase 14 replaces the static pin
; with the per-shape VRAM decision so the 24 GB shape reports cpu (gateway →
; gemini-stt) and >=30 GB shapes serve GPU whisper. WHISPER__MODEL_DIR /
; SPEACHES__MODEL_DIR stay removed — UAT 15 confirmed neither exists in WhisperConfig.
; Phase 06.8 STT fix iii (2026-05-21): the HF cache is now PRE-POPULATED by
; the pre-shaped weights tarball that onstart extracts to /weights/whisper
; (top-level models--Systran--faster-whisper-large-v3/{refs,snapshots/<hash>/}
; incl. README.md model card). HF_HUB_CACHE="/weights/whisper" points speaches'
; get_model_card_data_or_raise() gate at that cache — absolute path MUST match
; download-weights.sh's `tar -xzf … -C /weights/whisper` exactly (Pitfall 4). This
; supersedes the old "populated via POST /v1/models/{id}" runtime-fetch note.
; (The environment= line above carries HF_HUB_CACHE only; the whisper device
; is inherited from the onstart-exported env per the Phase 14 note above.)
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/speaches.log
priority=200

; The Infinity embed child was REMOVED — Phase 06.7 D-03: the embedding service
; relocated OFF the GPU pod to the 24/7 CPU host (n8n-ia-vm, GATE 4:
; multilingual-e5-large via Infinity). Phase 11.1 D-A2 also dropped the
; /opt/infinity-venv from the Dockerfile (rollback-only stage gone); embed
; stays off-pod through Phase 11.2.

[program:chatterbox]
; Phase 06.7 Plan 05: Chatterbox Multilingual TTS — 5th supervisord child,
; OpenAI /v1/audio/speech + GET /health on :8003 (gateway probe target).
; Isolated /opt/chatterbox-venv (torch 2.6.0 cu124 + transformers 5.2.0 +
; torchvision 0.21.0). directory=/opt/chatterbox-data so uvicorn imports
; chatterbox_server:app from the COPYed module. CHATTERBOX_* env is config-
; driven (paths/caps/prefix); MinIO creds reuse the pod's existing env for the
; zero-shot reference-WAV S3 fetch. priority=250 (after STT, before dcgm).
command=/opt/chatterbox-venv/bin/uvicorn chatterbox_server:app --host 0.0.0.0 --port 8003
directory=/opt/chatterbox-data
; HF_HUB_OFFLINE=1 + HF_HOME pinned to the cache dir: onstart pre-extracts the
; ResembleAI/chatterbox HF snapshot (models--ResembleAI--chatterbox/) from MinIO
; into CHATTERBOX_MODEL_CACHE_DIR, so from_pretrained()'s snapshot_download reads
; the local cache and never contacts huggingface.co. Before this (2026-06-13) the
; runtime HF fetch crash-looped on Vast hosts without an HF route, so /health on
; :8003 never came up and the gateway never marked the pod Ready.
environment=CHATTERBOX_MODEL_CACHE_DIR="/opt/chatterbox-data/models",CHATTERBOX_VOICE_CACHE_DIR="/opt/chatterbox-data/voices",CHATTERBOX_S3_VOICE_PREFIX="voices",CHATTERBOX_MAX_INPUT_CHARS="4000",CHATTERBOX_LANGUAGE_ID="pt",HF_HUB_OFFLINE="1",HF_HOME="/opt/chatterbox-data/models"
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/chatterbox.log
priority=250

[program:dcgm]
command=/usr/bin/dcgm-exporter --address 0.0.0.0:9400
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/dcgm.log
priority=400
