# pod/primary/supervisord.conf
; pod/primary/supervisord.conf — Phase 11.1 shrunk (STT off-pod).
;
; Phase 11.1 (2026-06-04) removed [program:speaches] (faster-whisper tier-0
; STT, now off-pod via tier-1 OpenAI fallback) per SEED-010 Mudança 2.
; The [program:infinity] commentary block that remained post-Phase 06.7 D-03
; is also gone (D-A2 bundle). Post-shrink child set: llama + chatterbox + dcgm
; (health-bridge runs separately in the pod, not under supervisord).
;
; Replaces the rejected DinD design (06.6-SPIKE-dind-privileged.md proved
; overlayfs in nested namespace fails on Vast.ai) with supervisord as
; PID 1 + co-located child processes (06.6-WAVE0-GATES.md Decision 2).
; Children share the container's network namespace (localhost ports),
; the GPU device, and the filesystem.
;
; The jinja flag on [program:llama] is the B1 GGUF-embedded path locked
; by 06.6-WAVE0-GATES.md Decision 3 + 06.6-SPIKE-qwen3.6-jinja.md Round 3
; (Qwen3.6 GGUF embedded chat_template renders tool-calling end-to-end
; via the jinja flag alone; no external template path needed).
;
; This file is baked into the image at build time (pod/primary/Dockerfile
; COPY directive). Runtime overrides are NOT supported in this plan —
; that would require regenerating supervisord.conf inside the pod and
; restarting supervisord, adding fragility. Deferred per CONTEXT.md
; Deferred Ideas to a future B2 fallback path.

[supervisord]
nodaemon=true
loglevel=info
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:llama]
command=/app/llama-server --host 0.0.0.0 --port 8000 -m /weights/qwen/model.gguf -ngl 99 -np 2 --ctx-size 16384 --jinja
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/llama.log
priority=100

[program:chatterbox]
; Phase 06.7 Plan 05: Chatterbox Multilingual TTS — 5th supervisord child,
; OpenAI /v1/audio/speech + GET /health on :8003 (gateway probe target).
; Isolated /opt/chatterbox-venv (torch 2.6.0 cu124 + transformers 5.2.0 +
; torchvision 0.21.0). directory=/opt/chatterbox-data so uvicorn imports
; chatterbox_server:app from the COPYed module. CHATTERBOX_* env is config-
; driven (paths/caps/prefix); MinIO creds reuse the pod's existing env for the
; zero-shot reference-WAV S3 fetch. priority=250 (after STT, before dcgm).
command=/opt/chatterbox-venv/bin/uvicorn chatterbox_server:app --host 0.0.0.0 --port 8003
directory=/opt/chatterbox-data
environment=CHATTERBOX_MODEL_CACHE_DIR="/opt/chatterbox-data/models",CHATTERBOX_VOICE_CACHE_DIR="/opt/chatterbox-data/voices",CHATTERBOX_S3_VOICE_PREFIX="voices",CHATTERBOX_MAX_INPUT_CHARS="4000",CHATTERBOX_LANGUAGE_ID="pt"
autorestart=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/chatterbox.log
priority=250

[program:dcgm]
command=/usr/bin/dcgm-exporter --address 0.0.0.0:9400
autorestart=true
stopwaitsecs=5
redirect_stderr=true
stdout_logfile=/var/log/dcgm.log
priority=400
